# root/etc/services.d/sleancache/run
#!/usr/bin/with-contenv bash

exec \
s6-setuidgid abc /acestrem.engine/start.sh
