# root/etc/services.d/aceproxy/run
#!/usr/bin/with-contenv bash

exec \
s6-setuidgid abc /HTTPAceProxy/acehttp.py
